# reinvent-2019/fully-automated-farm/device/pump/service_files/pigpio.service
[Unit]
Description = pigpio daemon

[Service]
ExecStart = /usr/bin/pigpiod
Type = forking
WorkingDirectory = /home/pi/agri-device/

[Install]
WantedBy = multi-user.target
