# reinvent-2019/lego-ev3-raspberry-pi-robot/docker/ev3dev/roscontroller.service
[Unit]
Description=EV3 ROS Controller
After=multi-user.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
WorkingDirectory=/home/robot/
User=robot
ExecStart=/bin/bash --login -c '/home/robot/service.sh'
